# scrapers/carparks/scrape.cron
* * * * * python /app/scrape.py $INFLUX_USER $INFLUX_PASSWORD
