# scrapers/river-levels/scrape.cron
* * * * * python /app/scrape.py $INFLUX_USER $INFLUX_PASSWORD > /proc/$(cat /var/run/crond.pid)/fd/1 2>&1
